# deployment/systemd/pi4-c922.service
[Unit]
Description=Edge Camera Pipeline (Pi4 + C922)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=RUST_LOG=info
ExecStart=/opt/edge-camera-pipeline/bin/pi4_c922_host --config /etc/edge-camera-pipeline/pi4-c922.toml
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

